# pi-fallback/services/scoreboard-fallback.service
[Unit]
Description=Highlands Latin School baseball scoreboard fallback server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__PI_USER__
WorkingDirectory=__REPO_ROOT__/pi-fallback
EnvironmentFile=-__REPO_ROOT__/pi-fallback/pi.env
ExecStart=__REPO_ROOT__/pi-fallback/.venv/bin/python __REPO_ROOT__/pi-fallback/app.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
