# services/scoreboard-wifi-failover.service
[Unit]
Description=Maintain preferred USB Wi-Fi with fallback for Highlands Latin School baseball scoreboard
After=NetworkManager.service network-online.target
Wants=NetworkManager.service

[Service]
Type=oneshot
WorkingDirectory=__APP_ROOT__
EnvironmentFile=-__APP_ROOT__/.env
ExecStart=/bin/bash __APP_ROOT__/scripts/maintain-wifi-failover.sh

