# services/scoreboard-local.service
[Unit]
Description=Highlands Latin School baseball scoreboard local server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_ROOT__
EnvironmentFile=-__APP_ROOT__/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=__APP_ROOT__/.venv/bin/python __APP_ROOT__/app.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
